# backup_structure/transclip/transclip.service
[Unit]
Description=TransClip - Speech to Text Transcription
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m transclip
Restart=on-failure
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority

[Install]
WantedBy=graphical-session.target 